# switchboard.service
[Unit]
Description=Switchboard MCP Server
After=network.target

[Service]
# Runs as root so it can spawn CC subprocesses as the restricted 'switchboard' user
# via the Agent SDK's user= parameter
WorkingDirectory=/opt/switchboard
ExecStart=/usr/bin/python3 -m switchboard
Environment=SWITCHBOARD_DB=/opt/switchboard/data/switchboard.db
Environment=SWITCHBOARD_PORT=8100
Environment=CLAUDE_CODE_STREAM_CLOSE_TIMEOUT=120000
# Secrets and site-specific config loaded from external file (not tracked in git).
# Create /etc/switchboard/env with: AUTH_ISSUER_URL, RESOURCE_URL, SLACK_BOT_TOKEN, SLACK_CHANNEL_ID
EnvironmentFile=/etc/switchboard/env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
