# switchboard.service
[Unit]
Description=Switchboard MCP Server
After=network.target

[Service]
# Runs as root so it can spawn CC subprocesses as the restricted 'switchboard' user
# via the Agent SDK's user= parameter
WorkingDirectory=/opt/switchboard
ExecStart=/usr/bin/python3 server.py
Environment=SWITCHBOARD_DB=/opt/switchboard/data/switchboard.db
Environment=SWITCHBOARD_PORT=8100
Environment=AUTH_ISSUER_URL=https://auth.stephenfritz.dev
Environment=RESOURCE_URL=https://switchboard.stephenfritz.dev/mcp
Environment=CLAUDE_CODE_STREAM_CLOSE_TIMEOUT=120000
Environment=SLACK_BOT_TOKEN=xoxb-30223900903-8408818373397-Y4yVL6DbwpTZ2ktsHrlBIOcM
Environment=SLACK_CHANNEL_ID=C0AKMQ2UY9L
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
